# trigger-build/trigger-build.ini
[program:trigger-build]
directory=/home/retune/retune-main/
command=/usr/bin/npm run service:trigger-build
environment=PORT="12345"
autostart=yes
autorestart=yes
